# etc/profile.d/nxaliases.sh
alias nxd='nxdisable'
alias nxds='nxdissite'
alias nxe='nxenable'
alias nxes='nxensite'
alias nxmk='nxmksite'
alias nxrl='nxreload'
alias nxrm='nxrmsite'
alias nxrs='nxrestart'
alias nxs='nxstart'
alias nxsp='nxstop'
alias nxst='nxstatus'
