# etc/profile.d/nxaliases.sh
alias nxdis='nxdissite'
alias nxens='nxensite'
alias nxmks='nxmksite'
alias nxrms='nxrmsite'
alias nxrl='nxreload'
alias nxrs='nxrestart'
alias nxs='nxstart'
alias nxst='nxstatus'
alias nxsp='nxstop'
