# etc/profile.d/nxpath.sh
export NXPATH="/etc/nginx"
export NXSITES="$NXPATH/sites-available"
export NXLOG="/var/log/nginx"
export NXERRORLOG="$NXLOG/error.log"
export NXACCESSLOG="$NXLOG/access.log"
